# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:arm_listener]
command=/bin/bash -c "/opt/arm/scripts/arm_listener.sh"

[program:armui]
command=python3 /opt/arm/arm/runui.py
